# scripts/agent-trigger-watcher.service
# OpenMates Agent Trigger Watcher — systemd user service
#
# Watches scripts/.agent-triggers/ for JSON trigger files written by the admin
# sidecar (Docker) and runs opencode plan-mode investigations on the host.
#
# Install:
#   mkdir -p ~/.config/systemd/user
#   cp scripts/agent-trigger-watcher.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now agent-trigger-watcher.service
#
# Check status:
#   systemctl --user status agent-trigger-watcher
#   journalctl --user -u agent-trigger-watcher -f
#
# Logs also written to: logs/agent-investigations.log

[Unit]
Description=OpenMates Agent Trigger Watcher
After=network.target

[Service]
Type=simple
ExecStart=/home/superdev/projects/OpenMates/scripts/agent-trigger-watcher.sh
Restart=always
RestartSec=10
WorkingDirectory=/home/superdev/projects/OpenMates
EnvironmentFile=-/home/superdev/projects/OpenMates/.env

[Install]
WantedBy=default.target
